# backend/crontab
# Cheetah SEPA cron — times are container-local; container TZ set to America/New_York.
# Min  Hour  DOM  Mon  DOW  Command
30     16    *    *    1-5  /usr/local/bin/python -m sepa.cli scan
30     8     *    *    1-5  /usr/local/bin/python -m sepa.cli brief
# Position-aware alerts — every 5 min during market hours (9:30–16:00 ET, Mon–Fri).
*/5    9-15  *    *    1-5  /usr/local/bin/python -m sepa.cli alerts
0,5,10,15,20,25,30 16 * *  1-5  /usr/local/bin/python -m sepa.cli alerts
