# launchd/com.cheetah.sepa.scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Cheetah SEPA — post-close nightly scan
  Runs at 5:00pm America/New_York Monday-Friday.

  Install:
    cp com.cheetah.sepa.scan.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.cheetah.sepa.scan.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.cheetah.sepa.scan.plist
    rm ~/Library/LaunchAgents/com.cheetah.sepa.scan.plist

  Note: launchd StartCalendarInterval uses local machine time — adjust if
  your laptop isn't on US/Eastern.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cheetah.sepa.scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ajay/clinet-test/cheetah-market-app/backend/.venv/bin/python</string>
        <string>-m</string>
        <string>sepa.cli</string>
        <string>scan</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ajay/clinet-test/cheetah-market-app/backend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/ajay/.cheetah/sepa-scan.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ajay/.cheetah/sepa-scan.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
